# roles/backend_upstart/templates/backend_cluster.conf.j2
start on [2345]
stop on [06]

script
    {{ REMOTE_HOME }}/{{ CLUSTER_APP_NAME }}/backend/target/universal/stage/bin/backend
end script