# roles/uwsgi/templates/git_commit_service.conf.j2
start on [2345]
stop on [06]

env MONGOLAB_URI={{ lookup('env', 'MONGOLAB_URI') }}

script
    cd {{ REMOTE_APP_HOME }}
    exec uwsgi --ini {{ CONFIG_DIR }}/git_commit_service.ini
end script