# roles/frontend_upstart/templates/frontend_cluster.conf.j2
start on [2345]
stop on [06]

script
    {{ REMOTE_HOME }}/{{ CLUSTER_APP_NAME }}/frontend/target/universal/stage/bin/frontend
end script